# local_ai_full_lab_v2/scripts/systemd/openwebui.service
[Unit]
Description=OpenWebUI (user)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
# Adjust ExecStart depending on how OpenWebUI was installed:
# - pip user install: %h/.local/bin/open-webui
# - pipx: %h/.local/pipx/venvs/open-webui/bin/open-webui
ExecStart=%h/.local/bin/open-webui serve --host 127.0.0.1 --port 3000
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
